# contrib/pks-openpgp-card.service
[Service]
ExecStart=/usr/bin/pks-openpgp-card
Environment=RUST_LOG=debug
